# build_files/supervisord.conf
[supervisord]
nodaemon=true

[program:redis]
command=/usr/local/bin/redis-server

[program:hornet]
command=/usr/bin/hornet